# generator/jobs/crontab
*/5 * * * * . ~/env/bin/activate && python ~/current/generator.py
