# generator/jobs/crontab
*/5 * * * * python ~/current/generator.py
