# generator/jobs/crontab
*/5 * * * * python ~/curent/generator.py
